# deployment/telegram-tori-bot.service
[Unit]
Description=Telegram Tori Bot
Documentation=https://github.com/raine/telegram-tori-bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tori
Group=tori
WorkingDirectory=/opt/telegram-tori-bot
ExecStart=/opt/telegram-tori-bot/start.sh
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/telegram-tori-bot
RestrictSUIDSGID=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
DevicePolicy=closed
RestrictRealtime=yes
LockPersonality=yes

# Logging (stdout/stderr go to journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-tori-bot

[Install]
WantedBy=multi-user.target
